# serial_gateway/serial_gateway.conf
# myservice - myservice job file

description "Serial Gateway"
author "James Moore <hello@jmoore.me>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

chdir /home/pi/Developer/hivemonitor/serial_gateway/

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn
respawn limit 10 5

# Run as the pi user
#setuid pi

# Run before process
#pre-start script
#  [ -d /var/run/myservice ] || mkdir -p /var/run/myservice
#  echo "Put bash code here"
#end script
        
# Start the process
exec /home/pi/Developer/hivemonitor/serial_gateway/serial_gateway.rb
